# odroid/vision-watcher.path
[Path]
PathModified=/root/vision.jar

[Install]
WantedBy=multi-user.target
